# systemd/iris-bot.service
[Unit]
Description=Iris Discord bot
After=network.target

[Service]
Type=simple
User=iris
Group=iris
WorkingDirectory=/opt/iris-bot
EnvironmentFile=/etc/iris-bot/iris-bot.env
ExecStart=/opt/iris-bot/.proto/bin/proto exec bun -- bun run /opt/iris-bot/src/index.ts
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
